# infrastructure/launchd_plist/com.investsys.train.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.investsys.train</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>-m</string>
        <string>training.train</string>
        <string>--config</string>
        <string>config/aws_config.json</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/investment-system</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/YOUR_USERNAME/investment-system</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/investment-system/logs/train.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/investment-system/logs/train_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
